# centos7/c7-systemd-postgresql11/etc/systemd/system/remain-after-exit.service
[Unit]
Description=RemainAfterExit service

[Service]
RemainAfterExit=true
ExecStart=/root/remain-after-exit.sh

[Install]
WantedBy=multi-user.target
